# deploy/supervisor/predict-arm-a.conf
[program:predict-arm-a]
command=/workspace/PREDICT/scripts/train_rl_instance.sh a
directory=/workspace/PREDICT
autostart=false
autorestart=unexpected
startsecs=30
stopasgroup=true
killasgroup=true
stopsignal=TERM
environment=CUDA_VISIBLE_DEVICES="0,1",HF_HOME="/workspace/.hf_home",WANDB_MODE="offline",PYTHONUNBUFFERED="1"
stdout_logfile=/workspace/PREDICT/outputs/arm_a_rl/supervisor.log
stdout_logfile_maxbytes=0
redirect_stderr=true
